# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:spark_master]
command=/opt/spark/start-spark.sh

[program:infinispan]
command=/opt/infinispan-server/bin/start-infinispan.sh



